# salt/gunicorn/gunicorn.conf
[program:gunicorn]
command=/var/www/toma/env/bin/python /var/www/toma/env/bin/gunicorn_django
directory=/var/www/toma
user=nobody
autostart=true
autorestart=true
redirect_stderr=True
